# aula-11-14/echo_uds_activated_service/echod_auto.socket
[Unit]
Description=Echo daemon activation socket
 

[Socket]
ListenStream=/opt/echod_auto/echod.sock
SocketMode=0660
SocketUser=root
SocketGroup=tvs

[Install]
WantedBy=sockets.target
