# aula-11-13/intro/b.service
[Unit] 
Description= test service b

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=sleep 1

[Install]
WantedBy=test.target
